# api-iam/devops/image/php-cli/consumers.conf
[program:create_user]
command=/var/www/app/bin/console rabbitmq:consumer create_user
autostart=true
autorestart=true
numprocs=10
process_name=%(program_name)s_%(process_num)02d